# cron/com.flatwhite.weekly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.flatwhite.weekly</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/TAC/Desktop/AntiGravity/FW Automate/flatwhite/cron/flatwhite_weekly.sh</string>
    </array>

    <!-- Wednesday 06:00 AEST = Tuesday 20:00 UTC -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>3</integer>
        <key>Hour</key>
        <integer>6</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/TAC/Desktop/AntiGravity/FW Automate/flatwhite/data/logs/cron.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/TAC/Desktop/AntiGravity/FW Automate/flatwhite/data/logs/cron.log</string>

    <key>WorkingDirectory</key>
    <string>/Users/TAC/Desktop/AntiGravity/FW Automate/flatwhite</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    </dict>
</dict>
</plist>
